# sample/kibana
#!/bin/bash 
### BEGIN INIT INFO
# Provides:          kibana
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Make sense of a mountain of logs.
### END INIT INFO

#
# chkconfig: 2345 20 80
#
KIBANA_PATH="/usr/share/kibana"

case "$1" in
start)
  ruby $KIBANA_PATH/kibana-daemon.rb $1
  RETVAL=$?
;;

stop)
  ruby $KIBANA_PATH/kibana-daemon.rb $1
  RETVAL=$?
;;

restart)
  ruby $KIBANA_PATH/kibana-daemon.rb $1
  RETVAL=$?
;;

status)
  ruby $KIBANA_PATH/kibana-daemon.rb $1
;;

force-reload)
  ruby $KIBANA_PATH/kibana-daemon.rb zap
  ruby $KIBANA_PATH/kibana-daemon.rb start
  RETVAL=$?
;;

*)
  echo "Usage: $0 {start|stop|restart|status|force-reload}"
  exit 1
esac

exit $RETVAL
